# debian/status-area-openvpn.postinst
#!/bin/sh

gtk-update-icon-cache -f /usr/share/icons/hicolor
rc-update del openvpn default 2>&1 > /dev/null || true

#DEBHELPER#
