# debian/status-area-openvpn.postrm
#!/bin/sh

update-rc.d openvpn enable || true

#DEBHELPER#
